# templates/apt-mirror.cron.j2
#!/bin/sh -e
# {{ ansible_managed }}

/usr/bin/apt-mirror
{{ apt_mirror_base_path }}/var/clean.sh
